# conf/etc/rc.local
#!/bin/sh -e

iptables-restore < /home/boj/iptables-save
exit 0
