# conf/etc/rc.local
#!/bin/sh -e

iptables-restore < /home/boj/iptables-save
[ -f /usr/local/sbin/pernode-rc.local ] && /usr/local/sbin/pernode-rc.local
exit 0
