# conf/etc/cron.d/ezvzdump
# backup at 3:20 every day
20 3 * * *  root   /home/boj/scripts/ezvzdump >/dev/null
